# .config/systemd/user/mpd-notify.service
[Unit]
Description=mpd-notify

[Service]
Type=simple
Restart=always
ExecStart=%h/.local/bin/mpd-notify

[Install]
WantedBy=default.target
